# catalyst/keykeeper/root_overlay/usr/local/keykeeper/distfiles/99_ldattach_rngd.rules
KERNEL=="ttyUSB*", RUN+="ldattach -8 -n -1 -s 921600 tty /dev/ttyUSB%n"
